# deploy/datastar-heap-live-map.service
[Unit]
Description=Datastar Heap View (live-map)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/datastar-hs
ExecStart=/opt/datastar-hs/heap-view live-map 4004
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
