# deploy/datastar-heap-live-fibs.service
[Unit]
Description=Datastar Heap View (live-fibs)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/datastar-hs
ExecStart=/opt/datastar-hs/heap-view live-fibs 4005
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
